# initscripts/devicejs
#! /bin/sh
# /etc/init.d/devicejs: starts DeviceJS & WigWag runtime

### BEGIN INIT INFO
# Provides:             devicejs
# Required-Start:       $remote_fs $time
# Required-Stop:        $remote_fs $time
# Should-Start:         $network $devicedbd
# Should-Stop:          $network $devicedbd
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    System logger
### END INIT INFO

# Source function library.
. /etc/init.d/functions

DEVICEJS_DIR="/wigwag/devicejs/web_applications"
DEVICEJS_CMD="./run-on-relay.sh webserver.js"

DEVICEJS_LOG="/wigwag/log/devicejs.log"

function congig_cpu_scaling() {
   MAX_FREQ=/sys/devices/system/cpu/cpu0/cpufreq/scaling_max_freq
   MIN_FREQ=/sys/devices/system/cpu/cpu0/cpufreq/scaling_min_freq
   cat $MAX_FREQ > $MIN_FREQ
}

function run_devjs() {
    if [ -e $DEVICEJS_LOG ]; then
	# some trivial log rotation
	mv $DEVICEJS_LOG $DEVICEJS_LOG.1
    fi
    
    pushd $DEVICEJS_DIR
    $DEVICEJS_CMD > $DEVICEJS_LOG 2>&1 &
#	$TUNSLIP -v$TUNSLIP_LOG_LEVEL -s $TUNSLIP_SERIAL $IP6_ROUTE > $TUNSLIP_LOG 2>&1 &
    popd
}

case "$1" in
    start)
	killall node
	eval $COLOR_BOLD
	echo "Starting DeviceJS services."
	eval $COLOR_NORMAL
        congig_cpu_scaling
	run_devjs
	;;
    stop)
	echo "Stopping DeviceJS"
	killall node
	;;
    restart) 
	killall node
        congig_cpu_scaling
	run_devjs
  	;;
      *)
	log_success_msg "Usage: $0 {start|stop|restart}"
    exit 1
esac

exit 0

     
