# initscripts/wwrelay
#! /bin/sh
# /etc/init.d/wwrelay: initialized the WigWag

### BEGIN INIT INFO
# Provides:             wwrelay
# Required-Start:       $remote_fs $time
# Required-Stop:        $remote_fs $time
# Should-Start:         $network
# Should-Stop:          $network
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    System logger
### END INIT INFO

# Source function library.
. /etc/init.d/functions

WIGWAGROOT="/wigwag"
GPIO_DIR="${WIGWAGROOT}/wwrelay-utils/GPIO"
WWRELAY_LOG="${WIGWAGROOT}/log/wwrelay.log"
GPIO_LOG="${WIGWAGROOT}/log/gpio.log"
DBDIR="/userdata/etc/devicejs/db"

function prepDB() {                          
  if [ ! -e ${DBDIR} ]; then
     mkdir -p ${DBDIR}
  fi                                                                    
}                                                                       

function genCerts() {

   if [ ! -e ${WIGWAGROOT}/devicejs-core-modules/AppServer/Certs ]; then
      mkdir ${WIGWAGROOT}/devicejs-core-modules/AppServer/Certs
   fi

   if [ ! -e ${WIGWAGROOT}/devicejs-core-modules/AppServer/Certs/ca.crt ]; then
      cd ${WIGWAGROOT}/devicejs-core-modules/AppServer/Certs
      pwgen 50 1 > passphrase
      openssl genrsa -des3 -passout file:passphrase -out ca.key 1024
      openssl req -new -passin file:passphrase -key ca.key -out ca.csr -subj "/C=US/ST=Texas/L=Austin/O=WigWag Inc./OU=Consumer/CN=www.wigwag.com/emailAddress=support@wigwag.com"
      openssl x509 -req -days 7300 -in ca.csr -out ca.crt -passin file:passphrase -signkey ca.key
      openssl genrsa -des3 -passout file:passphrase -out server.key 1024
      openssl req -new -passin file:passphrase -key server.key -out server.csr -subj "/C=US/ST=Texas/L=Austin/O=WigWag Inc./OU=Consumer/CN=www.wigwag.com/emailAddress=support@wigwag.com"
      cp server.key server.key.passphrase
      openssl rsa -in server.key.passphrase -passin file:passphrase -out server.key
      openssl x509 -req -days 7300 -in server.csr -passin file:passphrase -signkey server.key -out server.crt
      rm passphrase
   fi

}

case "$1" in
    start)
        /wigwag/wwrelay-utils/GPIO/led.sh 10 0 10
	echo "Starting WigWag Relay"
        prepDB
        genCerts
        if [ -x $WWRELAY_LOG ] ; then
            mv $WWRELAY_LOG $WWRLAY_LOG.1
        fi

        if [ -x $GPIO_LOG ] ; then
            mv $GPIO_LOG $GPIO_LOG.1
        fi

        cd ${WIGWAGROOT}/wwrelay-utils/I2C
        /usr/bin/node ./ww_eeprom_reader.js -c https://cloud.wigwag.com -o true -t ${WIGWAGROOT}/devicejs-core-modules/Runner/template.config.json -r ${WIGWAGROOT}/devicejs-core-modules/Runner/relay.config.json \
        -p ${WIGWAGROOT}/devicejs-core-modules/rsmi/radioProfile.template.json -s ${WIGWAGROOT}/devicejs-core-modules/rsmi/radioProfile.config.json > $WWRELAY_LOG 2>&1
        OUT=$?                                                                                
        if [ $OUT -eq 0 ]; then                                                                                                                                                           
          echo "EEPROM written properly, relay.config.json successfully generated"                               
        else                                                                                  
          echo "EEPROM not written properly, relay.config.json is not generated."                      
          /wigwag/wwrelay-utils/GPIO/led.sh 0 10 0                                      
        fi 
	;;	
    stop)
	echo "Stopping WigWag Relay"
esac

exit 0

     
